# aws/scripts/lifeset-backend.service
[Unit]
Description=LifeSet Backend API Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/lifeset-backend
Environment=NODE_ENV=production
EnvironmentFile=/opt/lifeset-backend/.env
ExecStart=/usr/bin/node /opt/lifeset-backend/dist/main.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/lifeset-backend/app.log
StandardError=append:/var/log/lifeset-backend/error.log

[Install]
WantedBy=multi-user.target
